# templates/etc/profile.d/pipeline.sh
export SPARK_NUM_SLAVES={{num_slaves}}
export SPARK_WORKER_MEMORY={{spark_worker_mem}}